# files/watchman.ini
[program:watchman]
command=/scripts/watchman.sh
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0